# systemd/hatpi-website.service
[Unit]
Description=HATPI Website (Flask via Gunicorn)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=hatuser
WorkingDirectory=/nfs/hatops/ar0/hatpi-website
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/gunicorn -c /nfs/hatops/ar0/hatpi-website/gunicorn.conf.py app:app
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target


